# cachet/crontab
* * * * * /usr/bin/php /srv/app/artisan schedule:run
