# cachet/crontab
* * * * * /usr/bin/php /srv/app/artisan schedule:run >> /dev/null 2>&
